# HWMonitorServer.socket
[Unit]
Description=HWMonitorServer socket

[Socket]
ListenStream=/run/HWMonitorServer.sock

[Install]
WantedBy=sockets.target